# audioCapture/etc/systemd/system/audioCapMonitor.service
[Unit]
Description=Publish audiocap health metrics (CPU temp, WiFi RSSI, load, uptime, capture status) via MQTT
After=network.target audioCap.service
Wants=audioCap.service

[Service]
Type=simple
User=jdn
Group=jdn
WorkingDirectory=/home/jdn/Code/AircraftAudioId/audioCapture
EnvironmentFile=-/etc/aircraftaudioid/mqtt.env
ExecStart=/home/jdn/Code/AircraftAudioId/aaiVenv/bin/python \
    /home/jdn/Code/AircraftAudioId/audioCapture/scripts/audioCapMonitor.py \
    --broker 192.168.166.5 \
    --interval 60 \
    --user ${MQTT_USER} \
    --password ${MQTT_PASSWORD}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
